# packages/bsp/jethub/05-jethub-set-mac.rules
SUBSYSTEM=="net",KERNEL=="eth0",RUN+="/usr/lib/armbian/jethub-set-mac %k"
